# crontab
* * * * * sudo /home/pi/symfoni-hockey/runSnow.sh
* * * * * /home/pi/symfoni-hockey/autopull.sh >> /home/pi/symfoni-hockey/cron.log

